# udev_rules_examples/99-mount-backup-disc.rules
ACTION=="add", SUBSYSTEM=="block", ATTRS{idVendor}=="0480", ATTRS{idProduct}=="a004", ENV{DEVTYPE}=="partition", RUN+="/opt/backup/run.sh"

